# script/vcptoolbox.service
[Unit]
Description=VCPToolBox Service
Documentation=https://github.com/lioensky/VCPToolBox
After=network.target

[Service]
Type=simple
User=zh
WorkingDirectory=/home/zh/projects/VCPToolBox

# 启动命令
ExecStart=/bin/bash -lc 'source "$HOME/.venv/bin/activate" && exec /home/zh/.nvm/versions/node/v22.15.1/bin/node /home/zh/projects/VCPToolBox/server.js'

# 自动重启
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
